# ops/backend.service.d/10-env.conf
[Service]
Environment=OPENAI_MODEL=gpt-4o-mini
# Environment=OPENAI_MODEL=gpt-4o-mini   # ← uncomment this instead if gpt-5 is not available
